# supervisor.ini
[program:gunicorn_two]
command=gunicorn -b 127.0.0.1:8000 --worker-class socketio.sgunicorn.GeventSocketIOWorker HalliganAvailability.wsgi
directory=/home/tyler/Documents/Dev_Work/halliganhelper
user=tyler
autostart=true
autorestart=true
redirect_stderr=true
autoreload=true

[supervisord]
identifier = supervisor

