# ansible/roles/halliganhelper/templates/supervisor.halliganhelper.ini
[group:halliganhelper]
environment =
    STATIC_ROOT="{{ static_root }}",
    MEDIA_ROOT="{{ media_root }}",
    EMAIL_PASSWORD="{{ email_password }}",
    DB_PASSWORD="{{ db_password }}",
    REDIS_PASSWORD="{{ redis_password }}",
    SECRET_KEY="{{ django_secret_key }}"

programs=hh-daphne,hh-workers

[program:hh-daphne]
user = hh
directory = /webapps/hh/src/halliganhelper
command = /webapps/hh/envs/current/bin/daphne halliganhelper.asgi:channel_layer
autostart = true
autorestart = true
stderr_logfile = /var/log/hh/daphne_stderr.log
stdout_logfile = /var/log/hh/daphne_stdout.log
stopsignal = INT

[program:hh-workers]
process_name = halliganhelper_runworker-%(process_num)s
user = hh
directory = /webapps/hh/src/halliganhelper
command = /webapps/hh/envs/current/bin/python manage.py runworker
numprocs = 2
autostart = true
autorestart = true
stderr_logfile = /var/log/hh/runworker_stderr.log
stdout_logfile = /var/log/hh/runworker_stdout.log
stopsignal = INT
